# etc/systemd/weather.timer
[Unit]
Description=Run FAA Weather Report

[Timer]
OnCalendar=*-*-* *:*/3:00

[Install]
WantedBy=timers.target