# etc/systemd/weather.timer
[Unit]
Description=Run FAA Weather Report

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target